# .config/autostart/,autostart.desktop
[Desktop Entry]
Type=Application
Name=,autostart
TryExec=,autostart
Exec=,autostart
NoDisplay=true
