# .config/autostart/,autostart.desktop
[Desktop Entry]
Exec=sh -c ,autostart
Name=,autostart
Terminal=False
Type=Application
